# renew-cert.timer
[Unit]
Description="Renew the cert"

[Timer]
OnBootSec=3min
OnUnitActiveSec=monthly
Unit=renew-cert.service

[Install]
WantedBy=multi-user.target