# renew-cert.service
[Unit]
Description="Renew the cert"

[Service]
ExecStart=/var/www/mikemaccana.com/renew-cert.bash