# renew-cert.timer
[Unit]
Description="Renew the cert"

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=multi-user.target